# packaging/systemd/mnemosd.service
[Unit]
Description=Mnemos Memory Daemon
Documentation=https://github.com/YourOrg/mnemos
# When socket-activated, systemd starts this unit on first connection.
Requires=mnemosd.socket
After=network-online.target

[Service]
Type=simple
# Look for the binary in common install locations.
ExecStart=%h/.cargo/bin/mnemosd
# Verify the daemon is serving before systemd considers the unit started.
ExecStartPost=/bin/sh -c 'for i in $(seq 1 30); do curl -sf http://127.0.0.1:7423/health >/dev/null && exit 0; sleep 1; done; exit 1'
Restart=on-failure
RestartSec=5
# Allow up to 3 minutes for model loading on slow machines.
TimeoutStartSec=180
# Limit resource usage.
MemoryMax=2G
# Logging.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mnemosd

[Install]
WantedBy=default.target
